# docker/yngdieng-mono/etc/services.d/gateway/run
#!/usr/bin/execlineb -P
/app/server/gateway/gateway_image.binary --port=8081 --endpoint=localhost:5002 --network=tcp