# docker/yngdieng-mono/etc/services.d/backend/run
#!/bin/sh
cd /app/server/backend/
dotnet Yngdieng.Backend.dll --IndexPath=/data/hotswap/ --TtsAudioFolder=/data/tts_audio_wav/
